# debian/postinst
#!/bin/sh
update-rc.d wb-prepare start 10 S

invoke-rc.d wb-prepare fix_macs

update-rc.d wb-gsm-rtc defaults
invoke-rc.d wb-gsm-rtc start

update-rc.d wb-init defaults
invoke-rc.d wb-init start

[ -e /var/www/uploads ] || mkdir -p /var/www/uploads
chown www-data:www-data /var/www/uploads
update-rc.d wb-watch-update defaults
invoke-rc.d wb-watch-update start

case "$1" in
install|upgrade|configure)
    wb_dir="/var/lib/wirenboard"
    short_sn_fname="$wb_dir/short_sn.conf"
    if [ ! -f "$short_sn_fname" ]; then
        echo "Storing short SN"
        if [ -e "/mnt/data/$short_sn_fname" ]; then
            cp "/mnt/data/$short_sn_fname" "$short_sn_fname"
        else
            short_sn=`wb-gen-serial -s`
            echo "$short_sn" > "$short_sn_fname"
        fi
    fi
esac