# debian/postinst
#!/bin/sh

rm -f "/var/run/wb_env.cache" || true

# The newer versions of init-system-helpers don't support
# update-rc.d start/stop and don't have insserv.
# In our case, this version info allows us to determine
# if package is being installed on wheezy or stretch.
# In wheezy we need to use update-rc.d start/stop,
# in stretch systemd will manage everything related to these services
# so we want to suppress all warnings related to update-rc.d.
INITHELPERS_VERSION=`dpkg-query -W -f='${Version}' init-system-helpers`
OLD_UPDATE_RC_D=true
if `dpkg --compare-versions $INITHELPERS_VERSION ge 1.22`; then
    OLD_UPDATE_RC_D=false
fi

if [ "$1" = "configure" ]; then
    if [ ! "$2" = "" ]; then
        if `dpkg --compare-versions $2 lt 1.75.0`; then
            # the runlevels for wb-prepare and wb-gsm were different
            # prior to 1.75
            # update-rc.d remove refuses to remove symlinks for wb-prepare
            # because of dependencies, so let's help it using low-level insserv
            if $OLD_UPDATE_RC_D; then
                insserv -f -r wb-gsm-rtc
            fi
        elif `dpkg --compare-versions $2 lt 1.73.1`; then
            # the runlevels for wb-gsm-rtc prior to 1.73.1 were different
            # this is necessary to clear the old symlinks. 
            # The new ones will be installed by update-rc.d below
            update-rc.d wb-gsm-rtc remove
        fi
    fi
fi

# FIXME: this is not in postinst style, but works without warning
# in stretch
setup_systemd_service() {
    systemctl enable $1 2>/dev/null
}

wb-prepare fix_macs

if $OLD_UPDATE_RC_D; then
    update-rc.d wb-gsm-rtc start 10 3 S . stop 10 0 6 . >/dev/null
else
    setup_systemd_service wb-gsm-rtc
fi

update-rc.d wb-init defaults >/dev/null
invoke-rc.d wb-init start

[ -e /var/www/uploads ] || mkdir -p /var/www/uploads
chown www-data:www-data /var/www/uploads
update-rc.d wb-watch-update defaults >/dev/null
invoke-rc.d wb-watch-update start

case "$1" in
install|upgrade|configure)
    wb_dir="/var/lib/wirenboard"
    short_sn_fname="$wb_dir/short_sn.conf"
    if [ ! -f "$short_sn_fname" ]; then
        echo "Reboot required for correct SN generation"
    fi
esac

#DEBHELPER#

exit 0
