# debian/postinst
#!/bin/sh
update-rc.d wb-prepare start 10 3 S
invoke-rc.d wb-prepare fix_macs

if [ "$1" = "configure" ]; then
    if [ ! "$2" = "" ]; then
        if `dpkg --compare-versions $2 lt 1.73.1`; then
            # the runlevels for wb-gsm-rtc prior to 1.73.1 were different
            # this is necessary to clear the old symlinks. 
            # The new ones will be installed by update-rc.d below
            update-rc.d wb-gsm-rtc remove
        fi
    fi
fi

update-rc.d wb-gsm-rtc start 10 3 S stop 10 0 6

update-rc.d wb-init defaults
invoke-rc.d wb-init start

[ -e /var/www/uploads ] || mkdir -p /var/www/uploads
chown www-data:www-data /var/www/uploads
update-rc.d wb-watch-update defaults
invoke-rc.d wb-watch-update start

case "$1" in
install|upgrade|configure)
    wb_dir="/var/lib/wirenboard"
    short_sn_fname="$wb_dir/short_sn.conf"
    if [ ! -f "$short_sn_fname" ]; then
        echo "Reboot required for correct SN generation"
    fi
esac

