# debian/wb-utils.wb-gsm.service
[Unit]
Description=Power on GSM modem
Wants=wb-hwconf-manager.service
Before=network-pre.target
After=wb-hwconf-manager.service

[Service]
Type=oneshot
ExecStart=wb-gsm mm_on
ExecStop=wb-gsm mm_off
ExecCondition=wb-gsm should_enable
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
